# ansible/roles/server/templates/kadeem-daemon.service.j2
[Unit]
Description=Kadeem API and Background Jobs
Documentation=https://github.com/galchammat/kadeem
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User={{ kadeem_user }}
Group={{ kadeem_group }}
WorkingDirectory={{ kadeem_home }}

# Environment variables
EnvironmentFile=/etc/kadeem/kadeem.env

# Start daemon
ExecStart={{ kadeem_bin_dir }}/daemon

# Restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kadeem-daemon

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ kadeem_home }}

[Install]
WantedBy=multi-user.target
